# IT490/APIServer.conf
#!upstart
description "API pull server"

start on startup
stop on shutdown

respawn

exec /home/amuthiyan/git/IT490/APIServer.php
